# systemd/pulse.service
[Unit]
Description=Security Pulse collector (system + threat feed snapshot)

[Service]
Type=oneshot
ExecStart=%h/.local/bin/pulse-collector
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7
